# files/Library/LaunchDaemons/com.brad-x.vmnet.vde2.plist
<?xml version='1.0' encoding='UTF-8'?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
"http://www.apple.com/DTDs/PropertyList-1.0.dtd" >
<plist version='1.0'>
	<dict>
		<key>Label</key>
		<string>com.brad-x.vmnet.vde2</string>
		<key>Program</key>
		<string>/opt/local/bin/vde_switch</string>
		<key>ProgramArguments</key>
		<array>
			<string>/opt/local/bin/vde_switch</string>
			<string>--pidfile</string>
			<string>/opt/local/var/run/vde2/vde_switch.pid</string>
			<string>--sock</string>
			<string>/opt/local/var/run/vde2/vde.ctl</string>
			<string>-M</string>
			<string>/opt/local/var/run/vde2/mgmt</string>
			<string>--nostdin</string>
		</array>
		<key>StandardErrorPath</key>
		<string>/opt/local/var/run/vde2/vde.stderr</string>
		<key>StandardOutPath</key>
		<string>/opt/local/var/run/vde2/vde.stdout</string>
		<key>Disabled</key><false/>
		<key>KeepAlive</key><true/>
	</dict>
</plist>