# files/Library/LaunchDaemons/com.brad-x.libvirt.virtlogd.plist
<?xml version='1.0' encoding='UTF-8'?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
"http://www.apple.com/DTDs/PropertyList-1.0.dtd" >
<plist version='1.0'>
	<dict>
		<key>Label</key>
		<string>com.brad-x.libvirt.virtlogd</string>
		<key>Program</key>
		<string>/opt/local/sbin/virtlogd</string>
		<key>ProgramArguments</key>
		<array>
			<string>/opt/local/sbin/virtlogd</string>
			<string>-f</string>
			<string>/opt/local/etc/libvirt/virtlogd.conf</string>
		</array>
		<key>StandardErrorPath</key>
		<string>/opt/local/var/run/virtlogd.stderr</string>
		<key>StandardOutPath</key>
		<string>/opt/local/var/run/virtlogd.stdout</string>
		<key>Disabled</key><false/>
		<key>KeepAlive</key><true/>
	</dict>
</plist>