# files/Library/LaunchDaemons/com.brad-x.vmnet.vde_vmnet.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- make install: yes -->
<plist version="1.0">
	<dict>
		<key>Label</key>
		<string>com.brad-x.vmnet.vde_vmnet</string>
		<key>Program</key>
		<string>/opt/local/bin/vde_vmnet</string>
		<key>ProgramArguments</key>
		<array>
			<string>/opt/local/bin/vde_vmnet</string>
			<string>--vmnet-gateway=192.168.127.1</string>
			<string>/opt/local/var/run/vde2/vde.ctl</string>
		</array>
		<key>StandardErrorPath</key>
		<string>/opt/local/var/run/vde2/vde_vmnet.stderr</string>
		<key>StandardOutPath</key>
		<string>/opt/local/var/run/vde2/vde_vmnet.stdout</string>
		<key>RunAtLoad</key>
		<true />
		<key>UserName</key>
		<string>root</string>
		<key>KeepAlive</key>
		<dict>
			<key>PathState</key>
			<dict>
				<key>/opt/local/var/run/vde2/vde_switch.pid</key>
				<true />
			</dict>
		</dict>
	</dict>
</plist>
